# ansible/roles/common/files/update-motd.service
[Unit]
Description=Update motd

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-motd.sh
User=root