# ansible/roles/borg_servers/files/compact_repo.timer
[Unit]
Description=Compact borg repos

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
