# ansible/roles/chaotic_aur/files/reflector-lower.timer
[Unit]
Description=Refresh Pacman mirrorlist weekly with Reflector.

[Timer]
OnCalendar=daily
Persistent=true
AccuracySec=1us
RandomizedDelaySec=12h

[Install]
WantedBy=timers.target